# systemd/beacon-alert-processor.service
[Unit]
Description=Beacon Alert Processor — NWS alert TTS conversion and cleanup
After=network.target
OnFailure=notify-email@%n.service

[Service]
Type=simple
WorkingDirectory=/home/lh_admin/weather_station
ExecStart=/home/lh_admin/weather_station/venv/bin/python run_alert_service.py
Restart=on-failure
RestartSec=10

Environment=MONGO_URI=mongodb://localhost:27017
Environment=ALERT_SERVICE_LOG=/home/lh_admin/weather_station/logs/alert_service.log

StandardOutput=append:/home/lh_admin/weather_station/logs/alert_service.log
StandardError=append:/home/lh_admin/weather_station/logs/alert_service.log

[Install]
WantedBy=default.target
