# systemd/fpren-bcp-ai-prefetch.service
[Unit]
Description=FPREN BCP AI Demographic Narrative Prefetch
After=network.target mongod.service

[Service]
Type=oneshot
User=ufuser
WorkingDirectory=/home/ufuser/Fpren-main
ExecStart=/home/ufuser/Fpren-main/venv/bin/python3 /home/ufuser/Fpren-main/weather_rss/bcp_ai_prefetch.py
Environment=MONGO_URI=mongodb://localhost:27017/
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
