# weather_rss/extended-fetcher.service
[Unit]
Description=FPREN Extended Data Fetcher (METAR, GA/NC Alerts, FL Traffic, School Closings)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ufuser/weather_rss
ExecStart=/home/ufuser/weather_rss/venv_rss/bin/python /home/ufuser/weather_rss/extended_fetcher.py
Restart=on-failure
RestartSec=10
Environment="LOG_FILE=/home/ufuser/weather_rss/logs/extended_fetcher_svc.log"

[Install]
WantedBy=default.target
